# supervisord.conf
[supervisord]
nodaemon=true

[program:apache]
command=/usr/local/bin/apache2-foreground
autorestart=true

[include]
files = /etc/supervisor/conf.d/*.conf

